# bin/rpi/dragonrise.service
[Unit]
Description=dragonrise systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /usr/local/sbin/dragonrise_service.sh

[Install]
WantedBy=multi-user.target
